# scripts/labrador.service
[Unit]
Description=Labrador Update
Before=systemd-user-sessions.service


[Service]
Type=idle
ExecStart=/usr/bin/python3 -u /home/pi/projects/labrador/scripts/gpio_functions.py

[Install]
WantedBy=multi-user.target
